# deploy/foil-records.timer
[Unit]
Description=Täglich Community-Rekorde snapshotten (Änderungen loggen + pushen)
[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
[Install]
WantedBy=timers.target
